# scripts/crontab
*/5 * * * * /mnt/config/cloudflare/update_ip_all
0 3 * * * /mnt/config/crontab/emby
